# packaging/konflikt.service
[Unit]
Description=Konflikt - Software KVM Switch
Documentation=https://github.com/Andersbakken/konflikt
After=graphical-session.target
Wants=avahi-daemon.service

[Service]
Type=simple
ExecStart=/usr/bin/konflikt
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=default.target
